# bin/systemd/style-web.target
[Unit]
StopWhenUnneeded=true
Wants=style-web-1.service
