# bin/systemd/style.target
[Unit]
StopWhenUnneeded=true
Wants=style-web.target

[Install]
WantedBy=multi-user.target
